# deploy/systemd/qr-warranty-worker.service
[Unit]
Description=24hStore QR Warranty Laravel Queue Worker
After=network.target mysql.service postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/24hstore-qr-warranty
ExecStart=/usr/bin/php artisan queue:work --sleep=3 --tries=3 --timeout=120 --max-time=3600
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=360

[Install]
WantedBy=multi-user.target
